# img/charon-sync.service
[Service]
Type=oneshot
ExecStart=/bin/bash /usr/charon_sync.sh
